# playbooks/templates/etc/init/kickstart.conf
description "Kickstart uWSGI app"
start on runlevel [2345]
stop on runlevel [06]

respawn

exec uwsgi --socket 127.0.0.1:3031 --plugin python --wsgi-file {{ kickstart_wsgi_app }} --master
